# aws/ami_configs/spot_setup.conf
description "start hostname setup"

# Don't expect long running process
task

start on (local-filesystems and net-device-up IFACE!=lo)
script
     /etc/spot_setup.sh
end script

