# ami_configs/ubuntu-12.04-x86_64-desktop-hvm/etc/rc.local
../../centos-65-x86_64-hvm-base/etc/rc.local